# root/etc/s6-overlay/s6-rc.d/x11vnc/run
#!/command/execlineb -P
x11vnc -display :99.0 --auth /tmp/xvfb-run -forever